# etc/profile.d/ros_python_logging.sh
export ROS_PYTHON_LOG_CONFIG_FILE="/etc/ros/python_logging.conf"
